# autorecommend/configs/stockrec-bot.service
[Unit]
Description=Stock Recommender - Telegram Bot
After=network.target

[Service]
Type=simple
User=$USER
WorkingDirectory=/home/$USER/stock-recommender/bot
Environment="PATH=/home/$USER/stock-recommender/venv/bin"
ExecStart=/home/$USER/stock-recommender/venv/bin/python src/bot_server.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
